# winalog-go/test_results/api_rules_01_list.json
{"rules":[{"id":"failed-login-threshold","name":"failed-login-threshold","description":"失败登录次数超过阈值","enabled":true,"severity":"high","score":70,"mitre_attack":["T1110"],"tags":["authentication","brute-force"],"is_custom":false,"event_ids":[4625],"levels":[2],"filter":{"event_ids":[4625],"levels":[2]},"message":"Failed login threshold exceeded: {{.Count}} failed attempts"},{"id":"admin-login-unusual","name":"admin-login-unusual","description":"管理员账户异常登录","enabled":true,"severity":"high","score":80,"mitre_attack":["T1078.004"],"tags":["authentication","privilege"],"is_custom":false,"event_ids":[4624],"levels":[4],"filter":{"event_ids":[4624],"levels":[4]},"message":"Unusual admin login detected for {{.User}} on {{.Computer}}"},{"id":"powershell-encoding","name":"powershell-encoding","description":"PowerShell编码命令执行","enabled":true,"severity":"high","score":85,"mitre_attack":["T1059.001"],"tags":["powershell","defense-evasion"],"is_custom":false,"event_ids":[4104],"levels":[4],"filter":{"event_ids":[4104],"levels":[4]},"message":"Encoded PowerShell command detected"},{"id":"security-log-cleared","name":"security-log-cleared","description":"安全日志被清除","enabled":true,"severity":"critical","score":100,"mitre_attack":["T1070.001"],"tags":["defense-evasion","clearing-logs"],"is_custom":false,"event_ids":[1102],"levels":[4],"filter":{"event_ids":[1102],"levels":[4]},"message":"Security log was cleared - potential attacker activity"},{"id":"account-lockout-alert","name":"account-lockout-alert","description":"账户被锁定","enabled":true,"severity":"medium","score":50,"mitre_attack":["T1078"],"tags":["authentication","account-management"],"is_custom":false,"event_ids":[4740],"levels":[4],"filter":{"event_ids":[4740],"levels":[4]},"message":"Account {{.User}} was locked out on {{.Computer}}"},{"id":"network-connection-alert","name":"network-connection-alert","description":"可疑网络连接","enabled":true,"severity":"medium","score":60,"mitre_attack":["T1059"],"tags":["network","command-and-control"],"is_custom":false,"event_ids":[3],"levels":[4],"filter":{"event_ids":[3],"levels":[4]},"message":"Suspicious network connection to port {{.DestinationPort}}"},{"id":"service-creation-alert","name":"service-creation-alert","description":"可疑服务创建","enabled":true,"severity":"high","score":75,"mitre_attack":["T1569.002"],"tags":["persistence","service"],"is_custom":false,"event_ids":[4697],"levels":[4],"filter":{"event_ids":[4697],"levels":[4]},"message":"Suspicious service created: {{.ServiceName}}"},{"id":"scheduled-task-creation","name":"scheduled-task-creation","description":"可疑计划任务创建","enabled":true,"severity":"high","score":75,"mitre_attack":["T1053.005"],"tags":["persistence","scheduled-task"],"is_custom":false,"event_ids":[4698],"levels":[4],"filter":{"event_ids":[4698],"levels":[4]},"message":"Suspicious scheduled task created: {{.TaskName}}"},{"id":"usb-device-insertion","name":"usb-device-insertion","description":"USB设备接入","enabled":true,"severity":"medium","score":50,"mitre_attack":["T1091"],"tags":["initial-access","physical"],"is_custom":false,"event_ids":[4663],"levels":[4],"filter":{"event_ids":[4663],"levels":[4]},"message":"USB device write detected on {{.Computer}}"},{"id":"pass-the-hash-suspect","name":"pass-the-hash-suspect","description":"Pass-the-Hash 可疑特征","enabled":true,"severity":"high","score":85,"mitre_attack":["T1550.002"],"tags":["credential-access","lateral-movement"],"is_custom":false,"event_ids":[4624,4672],"levels":[4],"filter":{"event_ids":[4624,4672],"levels":[4]},"message":"Possible Pass-the-Hash attack detected for {{.User}}"},{"id":"mimikatz-suspect","name":"mimikatz-suspect","description":"Mimikatz 相关工具执行","enabled":true,"severity":"critical","score":100,"mitre_attack":["T1003"],"tags":["credential-access","malware"],"is_custom":false,"event_ids":[1],"levels":[4],"filter":{"event_ids":[1],"levels":[4]},"message":"Mimikatz or credential dumping tool detected"},{"id":"encoded-powershell-alert","name":"encoded-powershell-alert","description":"PowerShell 远程编码命令","enabled":true,"severity":"high","score":85,"mitre_attack":["T1059.001"],"tags":["powershell","execution"],"is_custom":false,"event_ids":[4104],"levels":[4],"filter":{"event_ids":[4104],"levels":[4]},"message":"Encoded PowerShell remote command detected"},{"id":"remote-tasking-creation","name":"remote-tasking-creation","description":"远程任务创建 (SchTasks)","enabled":true,"severity":"medium","score":60,"mitre_attack":["T1053.005"],"tags":["persistence","remote-execution"],"is_custom":false,"event_ids":[4698],"levels":[4],"filter":{"event_ids":[4698],"levels":[4]},"message":"Remote scheduled task creation detected"},{"id":"suspicious-run-key","name":"suspicious-run-key","description":"可疑注册表Run键自启动","enabled":true,"severity":"high","score":80,"mitre_attack":["T1547.001"],"tags":["persistence","registry"],"is_custom":false,"event_ids":[4657],"levels":[4],"filter":{"event_ids":[4657],"levels":[4]},"message":"Suspicious Run key modification detected"},{"id":"suspicious-dll-load","name":"suspicious-dll-load","description":"可疑DLL加载","enabled":true,"severity":"medium","score":60,"mitre_attack":["T1574.001"],"tags":["defense-evasion","persistence"],"is_custom":false,"event_ids":[7],"levels":[4],"filter":{"event_ids":[7],"levels":[4]},"message":"Suspicious DLL loaded from unusual location"},{"id":"wmi-suspicious","name":"wmi-suspicious","description":"可疑WMI活动","enabled":true,"severity":"medium","score":65,"mitre_attack":["T1047"],"tags":["persistence","execution"],"is_custom":false,"event_ids":[5861],"levels":[4],"filter":{"event_ids":[5861],"levels":[4]},"message":"Suspicious WMI activity detected"},{"id":"registry-run-key-modification","name":"registry-run-key-modification","description":"注册表Run键自启动项被修改","enabled":true,"severity":"high","score":80,"mitre_attack":["T1547.001"],"tags":["persistence","registry"],"is_custom":false,"event_ids":[13],"levels":[4],"filter":{"event_ids":[13],"levels":[4]},"message":"Registry Run key modified: {{.TargetObject}}"},{"id":"registry-services-modification","name":"registry-services-modification","description":"注册表服务启动类型被修改","enabled":true,"severity":"high","score":75,"mitre_attack":["T1543.003"],"tags":["persistence","registry"],"is_custom":false,"event_ids":[13],"levels":[4],"filter":{"event_ids":[13],"levels":[4]},"message":"Service registry modified: {{.TargetObject}}"},{"id":"registry-hive-modification","name":"registry-hive-modification","description":"注册表Hive被修改(安全相关)","enabled":true,"severity":"critical","score":95,"mitre_attack":["T1003.002"],"tags":["credential-access","registry"],"is_custom":false,"event_ids":[13],"levels":[4],"filter":{"event_ids":[13],"levels":[4]},"message":"Critical registry hive modification: {{.TargetObject}}"},{"id":"registry-wmi-persistence","name":"registry-wmi-persistence","description":"WMI注册表持久化修改","enabled":true,"severity":"high","score":80,"mitre_attack":["T1546.003"],"tags":["persistence","wmi"],"is_custom":false,"event_ids":[13],"levels":[4],"filter":{"event_ids":[13],"levels":[4]},"message":"WMI registry persistence modification detected"}],"total_count":98,"enabled_count":98,"page":1,"page_size":20,"total_pages":5}